# mariadb-10.6-focal-auto-backups/etc/cron.d/dbbackup
MAILTO=""
#m      h       d	m	wd 	usr    	cmd
00      03      *	*	*	root	/opt/utils/backup-www-servers.sh