# mariadb-10.6-focal-auto-backups/etc/cron.d/dbbackup
MAILTO="" 
#m	h	d	m	wd	usr     cmd
0   3   *   *   *	root    /opt/dbbackup/do-backup.sh