# mariadb-11-noble-auto-backups/etc/cron.d/dbbackup
MAILTO=''
#m      h	d	m	wd 	usr    	cmd
30	3	*	*	*	root	/opt/dbbackup/do-backup.sh
